# deploy/systemd/streetscan-api.service
[Unit]
Description=Buqata StreetScan FastAPI
After=network.target postgresql.service

[Service]
User=streetscan
Group=streetscan
WorkingDirectory=/opt/buqata-streetscan/backend
EnvironmentFile=/opt/buqata-streetscan/backend/.env
ExecStart=/opt/buqata-streetscan/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8005
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
