# deploy/systemd/streetscan-worker.service
[Unit]
Description=Buqata StreetScan YOLO worker
After=network.target postgresql.service

[Service]
User=streetscan
Group=streetscan
WorkingDirectory=/opt/buqata-streetscan/backend
EnvironmentFile=/opt/buqata-streetscan/backend/.env
ExecStart=/opt/buqata-streetscan/backend/.venv/bin/python -m app.worker
Restart=always
RestartSec=10
# YOLO inference on CPU can be hungry; keep it polite next to the other apps
# on this shared 4GB box. MemoryMax kills+restarts the worker instead of
# letting it OOM a neighbour service.
Nice=10
CPUWeight=50
MemoryHigh=1200M
MemoryMax=1536M

[Install]
WantedBy=multi-user.target
